# release/fs/overlay/etc/profile.d/SSEnv.sh
export HOME=/
export PATH=$PATH:/oem/usr/bin:/system/bin
export LD_LIBRARY_PATH=/oem/usr/lib
