# release/fs/overlay/etc/profile.d/SSEnv.sh
export HOME=/userdata
export PATH=$PATH:/system/bin/:/oem/bin
export LD_LIBRARY_PATH=/oem/lib:/system/lib
